# systemd/domoticz.service.d/restart.conf
[Service]
Type=forking
PIDFile=/run/domoticz.pid
RemainAfterExit=no
Restart=on-failure
RestartSec=5s
